# deployment/systemd/youtube-scraper.service
[Unit]
Description=YouTube Scraper Production Service
Documentation=https://github.com/canaanhowell/youtube-scraper-production
After=network.target docker.service
Wants=network.target
Requires=docker.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/youtube_app

# Environment
Environment=PYTHONPATH=/opt/youtube_app
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/youtube_app/.env

# Script execution
ExecStart=/opt/youtube_app/deployment/run_scraper.sh

# Process management
RemainAfterExit=no
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=1800
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=false

# Logging
StandardOutput=append:/opt/youtube_app/logs/scraper-service.log
StandardError=append:/opt/youtube_app/logs/scraper-service-error.log
SyslogIdentifier=youtube-scraper

# Health monitoring
WatchdogSec=300
NotifyAccess=none

[Install]
WantedBy=multi-user.target