# deployment/systemd/youtube-scraper.timer
[Unit]
Description=Run YouTube Scraper every hour
Documentation=https://github.com/canaanhowell/youtube-scraper-production
Requires=youtube-scraper.service

[Timer]
# Run every hour at 15 minutes past the hour
OnCalendar=*:15
# Run immediately on timer activation if missed
Persistent=true
# Randomize by up to 5 minutes to avoid exact timing
RandomizedDelaySec=300

[Install]
WantedBy=timers.target