# deployment/systemd/youtube-analytics-fixed.service
[Unit]
Description=YouTube Analytics Processing Service
Documentation=https://github.com/canaanhowell/youtube-scraper-production
After=network.target youtube-scraper.service
Wants=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/youtube_app

# Environment
Environment=PYTHONPATH=/opt/youtube_app
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/youtube_app/.env

# Simple execution without problematic ExecStartPre
ExecStart=/bin/bash -c 'cd /opt/youtube_app && source venv/bin/activate && python -m src.scripts.collectors.run_analytics --task all'

# Process management
TimeoutStartSec=1800
TimeoutStopSec=60
RemainAfterExit=no

# Resource limits
LimitNOFILE=65536
LimitNPROC=2048

# Logging
StandardOutput=append:/opt/youtube_app/logs/analytics-service.log
StandardError=append:/opt/youtube_app/logs/analytics-service-error.log
SyslogIdentifier=youtube-analytics

[Install]
WantedBy=multi-user.target