# deployment/systemd/youtube-monitoring.timer
[Unit]
Description=YouTube Scraper Health Monitoring Timer
Documentation=https://github.com/canaanhowell/youtube-scraper-production
Requires=youtube-monitoring.service

[Timer]
# Run health checks every 30 minutes
OnCalendar=*:0/30:00
# Run 5 minutes after system boot
OnBootSec=5min
# If the system was down when the timer should have run, run it on next boot
Persistent=true
# Small randomization to avoid exact timing conflicts
RandomizedDelaySec=60

[Install]
WantedBy=timers.target